# config/systemd/vitess-cluster.target
[Unit]
Description=Vitess cluster Mode
Requires=network.service